# deploy/com.rarebook.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rarebook.web</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/chiwon/rare-book-store/venv/bin/gunicorn</string>
    <string>--bind</string>
    <string>127.0.0.1:8000</string>
    <string>--workers</string>
    <string>2</string>
    <string>--timeout</string>
    <string>120</string>
    <string>run:app</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/chiwon/rare-book-store</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/chiwon/rare-book-store/logs/gunicorn.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chiwon/rare-book-store/logs/gunicorn.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
